# restserver.conf
# restserver service script
start on filesystem
exec python3  /home/vagrant/restserver/src/server.py
